# Docker_Tailscale/alpine-ethtool-configuration
#!/sbin/openrc-run
# Stored in /etc/init.d/ethtool-configuration
# curl -L -O https://raw.githubusercontent.com/5points/Docker_tools_Script/main/Docker_Tailscale/alpine-ethtool-configuration /etc/init.d/ethtool-configuration
# chmod +x /etc/init.d/ethtool-configuration
# rc-update update add ethtool-configuration && /etc/init.d/ethtool-configuration start

COMMAND_PATH="/etc/networkd-dispatcher/routable.d/50-tailscale"

name="Tailscaled ethtool-configuration"
description="Runs the 50-tailscale script"
command_background="no"

depend() {
    after networking
    need net
    keyword="oneshot"
}

start() {
    ebegin "Starting Tailscaled ethtool-configuration"
    /etc/networkd-dispatcher/routable.d/50-tailscale
    eend $?
}

stop() {
    # Do nothing, this is a one-time script
    :
}
